# deploy/systemd/chordstash.service
[Unit]
Description=ChordStash Chord Sheet Manager
After=network.target

[Service]
Type=simple
User=administrator
WorkingDirectory=/home/administrator/AI/ChordStash
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=9010

[Install]
WantedBy=multi-user.target
